# vps/systemd/voice-sidecar.service
[Unit]
Description=Voice Sidecar (transcript echo + delete)
After=openclaw-gateway.service

[Service]
ExecStart=/usr/bin/python3 /root/voice-sidecar/main.py
Restart=always
RestartSec=5
Environment=HOME=/root
EnvironmentFile=-/root/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
